# .sisyphus/deploy/ai-frontend.service
[Unit]
Description=AI Practice Frontend Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-practice/web
Environment=NODE_ENV=production
# single-node deploy contract: one loopback frontend instance on 127.0.0.1:3445 behind nginx.
# future multi-instance rollout still needs external ingress/LB orchestration; this unit is node-local only.
ExecStart=/usr/bin/npm run start -- -p 3445 -H 127.0.0.1
Restart=always
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
